# packaging/debian/postinst
#!/bin/bash
set -e

AGORA_USER="agora"
AGORA_BASE="/opt/agora"

# ── Create user if needed ──
if ! id "${AGORA_USER}" &>/dev/null; then
    useradd -r -m -s /bin/bash "${AGORA_USER}"
fi

# ── Create persist directory (flash-backed, survives reboots) ──
mkdir -p "${AGORA_BASE}/persist"
chown "${AGORA_USER}:${AGORA_USER}" "${AGORA_BASE}/persist"

# ── Migrate persistent files from old state dir if needed ──
for file in cms_auth_token cms_config.json splash api_key; do
    src="${AGORA_BASE}/state/${file}"
    dst="${AGORA_BASE}/persist/${file}"
    if [ -f "${src}" ] && [ ! -f "${dst}" ]; then
        cp -p "${src}" "${dst}"
        echo "Agora: migrated ${file} from state/ to persist/"
    fi
done

# ── Mount state directory as tmpfs (RAM-backed, prevents flash wear) ──
FSTAB_LINE="tmpfs ${AGORA_BASE}/state tmpfs nodev,nosuid,noexec,size=10m,uid=${AGORA_USER},gid=${AGORA_USER},mode=0755 0 0"
if ! grep -q "${AGORA_BASE}/state" /etc/fstab 2>/dev/null; then
    echo "${FSTAB_LINE}" >> /etc/fstab
    echo "Agora: added tmpfs mount for ${AGORA_BASE}/state"
fi
mkdir -p "${AGORA_BASE}/state"
if ! mountpoint -q "${AGORA_BASE}/state"; then
    mount "${AGORA_BASE}/state"
fi

# ── Ensure directory ownership ──
chown -R "${AGORA_USER}:${AGORA_USER}" "${AGORA_BASE}/assets"
chown -R "${AGORA_USER}:${AGORA_USER}" "${AGORA_BASE}/state"
chown -R "${AGORA_USER}:${AGORA_USER}" "${AGORA_BASE}/logs"
chown -R "${AGORA_USER}:${AGORA_USER}" "${AGORA_BASE}/tmp"

# ── Generate config if first install ──
BOOT_CONFIG="/boot/agora-config.json"
if [ ! -f "${BOOT_CONFIG}" ]; then
    API_KEY=$(python3 -c "import secrets; print(secrets.token_urlsafe(32))")
    SECRET_KEY=$(python3 -c "import secrets; print(secrets.token_urlsafe(32))")
    cat > "${BOOT_CONFIG}" <<EOF
{
    "api_key": "${API_KEY}",
    "web_username": "admin",
    "web_password": "agora",
    "secret_key": "${SECRET_KEY}",
    "device_name": "$(hostname)",
    "default_splash": "splash/default.png"
}
EOF
    chmod 600 "${BOOT_CONFIG}"
    echo "Agora: generated config at ${BOOT_CONFIG} (default password: agora)"
fi

# ── Install Python packages ──
pip3 install --break-system-packages -q \
    -r "${AGORA_BASE}/src/requirements-api.txt" \
    -r "${AGORA_BASE}/src/requirements-player.txt" \
    -r "${AGORA_BASE}/src/requirements-cms-client.txt" 2>/dev/null \
|| pip3 install -q \
    -r "${AGORA_BASE}/src/requirements-api.txt" \
    -r "${AGORA_BASE}/src/requirements-player.txt" \
    -r "${AGORA_BASE}/src/requirements-cms-client.txt"

# ── Plymouth boot splash ──
if command -v plymouth-set-default-theme &>/dev/null; then
    plymouth-set-default-theme splash 2>/dev/null || true
fi

# ── Kernel command line: quiet boot + splash ──
CMDLINE="/boot/firmware/cmdline.txt"
if [ -f "${CMDLINE}" ]; then
    changed=false
    for flag in quiet splash plymouth.ignore-serial-consoles logo.nologo "vt.global_cursor_default=0"; do
        if ! grep -q "${flag}" "${CMDLINE}"; then
            sed -i "s/$/ ${flag}/" "${CMDLINE}"
            changed=true
        fi
    done
    if $changed; then
        echo "Agora: patched ${CMDLINE} for quiet boot + splash"
    fi
fi

# ── Boot config: HDMI force + overscan ──
BOOT_CFG="/boot/firmware/config.txt"
if [ -f "${BOOT_CFG}" ]; then
    for setting in disable_fw_kms_setup=1 disable_overscan=1 hdmi_force_hotplug=1 hdmi_force_edid_audio=1; do
        key="${setting%%=*}"
        if ! grep -q "^${key}=" "${BOOT_CFG}"; then
            echo "${setting}" >> "${BOOT_CFG}"
            echo "Agora: added ${setting} to ${BOOT_CFG}"
        fi
    done
fi

# ── Disable console login on HDMI (tty1) ──
systemctl disable getty@tty1 2>/dev/null || true

# ── Rebuild initramfs with new splash theme ──
update-initramfs -u 2>/dev/null || true

# ── Enable and restart services ──
systemctl daemon-reload
systemctl enable agora-player agora-api agora-cms-client
if [ "$1" = "configure" ]; then
    systemctl restart agora-player || true
    systemctl restart agora-api || true
    systemctl restart agora-cms-client || true
fi
