# systemd/agora-api.service
[Unit]
Description=Agora API Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=mpv
ExecStart=/usr/bin/python3 -m uvicorn api.main:app --host 0.0.0.0 --port 8000 --workers 1 --timeout-keep-alive 300
WorkingDirectory=/opt/agora/src
Environment=PYTHONPATH=/opt/agora/src
Environment=AGORA_BASE=/opt/agora
Environment=TMPDIR=/opt/agora/tmp
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
