# systemd/agora-provision.service
[Unit]
Description=Agora Provisioning (Captive Portal)
After=NetworkManager.service
Wants=NetworkManager.service
Before=agora-api.service agora-cms-client.service

[Service]
Type=simple
ExecStart=/bin/bash -c 'exec /usr/bin/python3 -u /opt/agora/src/provision/service.py 2>&1 | tee -a /boot/firmware/provision.log'
WorkingDirectory=/opt/agora/src
Environment=PYTHONPATH=/opt/agora/src
Environment=AGORA_BASE=/opt/agora
ExecStartPre=/bin/bash -c 'echo "=== provision start $(date) ===" >> /boot/firmware/provision.log'
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Needs root for port 80 and NetworkManager control
User=root

[Install]
WantedBy=multi-user.target
