# systemd/agora-player.service
[Unit]
Description=Agora Media Player
After=opt-agora-state.mount
Requires=opt-agora-state.mount
DefaultDependencies=no

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/agora/src/player/main.py
WorkingDirectory=/opt/agora/src
Environment=PYTHONPATH=/opt/agora/src
Environment=AGORA_BASE=/opt/agora
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
