# systemd/agora-cms-client.service
[Unit]
Description=Agora CMS Client
After=network-online.target agora-api.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 -m cms_client.main
WorkingDirectory=/opt/agora/src
Environment=PYTHONPATH=/opt/agora/src
Environment=AGORA_BASE=/opt/agora
Restart=always
RestartSec=10

# Configure NTP to use CMS server, then wait for network to settle
ExecStartPre=/bin/bash /opt/agora/src/scripts/configure-ntp.sh
ExecStartPre=/bin/sleep 5

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
